# dotfiles/autostart/conky_new.desktop
[Desktop Entry]
Name=Conky
Comment=Conky overlays
Exec=sh -c "sleep 5 && conky -c ~/.config/conky/conky_new.conf"
Terminal=false
Type=Application