# dotfiles/autostart/conky.desktop
[Desktop Entry]
Name=Conky
Comment=Conky overlays
Exec=sh -c "~/.config/conky/run_conky.sh"
Terminal=false
Type=Application