# dotfiles/autostart/conky_upcoming.desktop
[Desktop Entry]
Name=Conky
Comment=Conky overlays
Exec=sh -c "sleep 5 && conky -c ~/.config/conky/conky_upcoming.conf"
Terminal=false
Type=Application