# dotfiles/autostart/steam-charts.desktop
[Desktop Entry]
Name=Steam charts
Comment=Steam charts python script
Exec=sh -c "~/.config/conky/run_steam_script.sh"
Terminal=false
Type=Application